# services.d/reflection/run
#!/usr/bin/with-contenv bash
WEBUI_PORT=${WEBUI_PORT:-8080}

exec \
s6-setuidgid abc /tmp/main --api-addr="http://localhost:${WEBUI_PORT}/"
