# services.d/reflection/run
#!/usr/bin/with-contenv bash

exec \
s6-setuidgid abc /tmp/Reflection/main --api-addr="http://localhost:8088/"
